# supervisord.conf
[program:insta-js]
command=npm start
directory=/app
process_name=insta-js
stdout_logfile=/app/log
stderr_logfile=/app/log